# host/ai-fortress-key-sweep.service
[Unit]
Description=AI Fortress orphan virtual-key sweep
After=ai-fortress-litellm.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fortress-sweep
